# supervisord.conf
[supervisord]
nodaemon=true

[program:laravel]
command=php artisan serve --host=0.0.0.0 --port=8000
directory=/var/www
autostart=true
autorestart=true
stdout_logfile=/var/www/storage/logs/laravel.log
stderr_logfile=/var/www/storage/logs/laravel_error.log

[eventlistener:laravel-reloader]
command=sh -c 'while inotifywait -e modify,create,delete,move /var/www; do supervisorctl restart laravel; done'
events=PROCESS_STATE_EXIT
